# templates/desktop-apps/desktop-app.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=desktop app
PartOf=graphical-session.target

[Service]
ExecStart=%i

[Install]
WantedBy=graphical-session.target
